# etc/supervisor/conf.d/nginx.conf
[program:nginx]
command=/usr/local/sbin/nginx -c /etc/nginx/nginx.conf -g "daemon off;"
autostart=true
autorestart=true
